# data/indicator-keyboard.service.in
[Unit]
Description=Indicator Keyboard Backend
PartOf=graphical-session.target
PartOf=indicators.target

[Service]
ExecStart=@pkglibexecdir@/indicator-keyboard-service --use-gtk
Restart=on-failure

[Install]
WantedBy=indicators.target
